# scripts/com.ki-tracker.handoff-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ki-tracker.handoff-check</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Library/WebServer/Documents/KI Usage tracker/scripts/check-handoff.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/ki-tracker-handoff.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/ki-tracker-handoff.log</string>
</dict>
</plist>
